# update.service
[Unit]
Description=Startup

[Service]
ExecStart=/usr/local/sbin/update

[Install]
WantedBy=multi-user.target
